# elasticsearch/supervisor/eswrapper.ini
[supervisord]
user=root
nodaemon=true

[program:eswrapper]
user=root
command=docker-entrypoint.sh eswrapper
