# elasticsearch/supervisor/eswrapper.ini
[supervisord]
nodaemon=true

[program:eswrapper]
command=eswrapper
