# kinaba/supervisor/kibana.ini
[supervisord]
user=root
nodaemon=true

[program:kibana]
user=kibana
command=/usr/local/bin/kibana-docker