# python/amd64/supervisor/sshd.ini
[supervisord]
user=root
nodaemon=true

[program:sshd]
user=root
command=/usr/sbin/sshd -D