# mongo/supervisor/mongo.ini
[supervisord]
user=root
nodaemon=true

[program:mongo]
user=root
command=/usr/local/bin/docker-entrypoint.sh mongod
