# openresty/supervisor/openresty.ini
[supervisord]
nodaemon=true

[program:openresty]
command=/usr/bin/openresty -g 'daemon off;'
