# mariadb/supervisor/mariadb.ini
[supervisord]
nodaemon=true

[program:mariadb]
command=docker-entrypoint.sh mysqld
priority=10
autorestart=true
startretries=3